# software/system/src/lcdmon/lcdmonitor.upstart.conf
# lcdmonitor - lcdmonitor
#
# This service keeps lcdmonitor running

start on runlevel [345]
stop on runlevel [0126]
respawn
exec /usr/local/sbin/lcdmonitor  >/dev/null 2>&1