# software/system/src/ppsd/ppsd.upstart.conf
# ppsd - ppsd
#
# This service keeps ppsd running

start on runlevel [345]
stop on runlevel [0126]
respawn
exec /usr/local/sbin/ppsd >/dev/null 2>&1
